# services/gattrose-orchestrator.service
[Unit]
Description=Gattrose-NG Orchestrator Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/gattrose-ng
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/gattrose-ng/.venv/bin/python /opt/gattrose-ng/src/services/orchestrator.py
Restart=always
RestartSec=10

# Security
PrivateTmp=true
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/gattrose-ng/data /tmp

# Capabilities needed for wireless operations
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
